# scripts/nightly/deploy/launchd/com.luthien.nightly.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS LaunchAgent for luthien-proxy nightly maintenance.
     Rendered by deploy/install.sh. Placeholders:
       __LABEL__         com.luthien.nightly
       __NIGHTLY_SH__    absolute path to nightly.sh
       __HOUR__          0-23
       __MINUTE__        0-59
       __WORKING_DIR__   absolute path to repo root (for cwd)
       __LOG_DIR__       absolute path for stdout/stderr logs
       __PATH__          PATH to use (so brew/uv/claude/gh are visible)
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>__NIGHTLY_SH__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__WORKING_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>__HOUR__</integer>
    <key>Minute</key><integer>__MINUTE__</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/nightly.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/nightly.err.log</string>
</dict>
</plist>
